# servers/cron/tornado.textpattern.net/www-data.cron
15,45 * * * * bash /var/www/servers/textpattern.com/design-patterns/staging/build.sh
