# servers/cron/tribune.textpattern.com/www-data.cron
15,45 * * * * bash /var/www/sites/textpattern.com/design-patterns/staging/build.sh
