# servers/cron/tractor.textpattern.com/www-data.cron
0,15,30,45 * * * * bash /var/www/vhosts/textpattern.com/design-patterns/staging/build.sh
