# servers/cron/tractor.textpattern.com/www-data.cron
0,15,30,45 * * * * bash /var/www/vhosts/textpattern.com/design-patterns/staging/build.sh
0,15,30,45 * * * * sleep 5 && bash /var/www/vhosts/textpattern.com/design-patterns-4-6/staging/build.sh
0,15,30,45 * * * * sleep 10 && bash /var/www/vhosts/textpattern.com/design-patterns-4-7/staging/build.sh
