# live/servers/cron/tarzan.textpattern.net/www-data.cron
0,15,30,45 * * * * bash /var/www/servers/textpattern.com/docs-test/staging/mv-perms.sh
